# ha_docker/rootfs/etc/services.d/app/finish
#!/usr/bin/with-contenv bash
echo "[RF ADDON] Service stopped."
sleep 1
